# config/systemd/user/dunst.service.d/override.conf
[Unit]
PartOf=autostart.target

[Install]
WantedBy=
WantedBy=autostart.target
