# config/systemd/user/backup-home.service
[Unit]
Description=Backup $HOME directory

[Service]
Type=simple
ExecStart=%h/.local/bin/backup --quiet
