# config/systemd/user/backup-home.timer
[Unit]
Description=Backup $HOME directory hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
