# config/systemd/user/dex.service
[Unit]
Description=Autostart applications

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=dex --autostart
