# tag-lnx/config/systemd/user/backup-home.timer
[Unit]
Description=Backup home directory hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
